# airootfs/aai/db/modules/etc/udev/rules.d/96-fb-all-vcs-setup.rules
# Setup all vconsoles for a new framebuffer device
KERNEL=="fb*", ACTION=="add", RUN+="/etc/udev/all-vcs-set.sh"
