# root-image/aai/db/modules/etc/rc.local
#!/bin/sh
for TTY in /dev/tty?
do
	/usr/bin/setleds -D +num < "${TTY}"
done
