# root-image/aai/db/modules/usr/local/lib/systemd/user/jack.service
[Unit]
Description=JACK
After=sound.target
[Service]
ExecStart=/usr/bin/dbus-launch /usr/bin/jackd -R -dalsa
[Install]
WantedBy=default.target
