# root-image/aai/db/modules/usr/local/lib/systemd/user/timidity.service
[Unit]
Description=TiMidity
After=sound.target
[Service]
ExecStart=/usr/bin/timidity -iA -OO
[Install]
WantedBy=default.target
